# opi_web.service
[Unit]
Description=LoxIO Core Web UI
Documentation=https://github.com/Azazel101/orangepi-zero3-gpio-api
After=network.target opi_gpio.service
Wants=opi_gpio.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/opi_gpio_app/web

# Startup
ExecStart=/usr/bin/python3 app.py

# Restart policy with backoff
Restart=on-failure
RestartSec=5
RestartMaxDelaySec=60

# Timeout protection
TimeoutStartSec=20
TimeoutStopSec=10
KillMode=mixed
KillSignal=SIGTERM

# Resource limits
MemoryMax=256M
TasksMax=128

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=loxio-web

[Install]
WantedBy=multi-user.target
